# .config/systemd/user/copyq.service
[Unit]
Description=copyq (clipboard manager)

[Service]
ExecStart=copyq
Restart=always

[Install]
WantedBy=user-graphical-login.target
