# .config/systemd/user/hydroxide-smtp.service
[Unit]
Description=hydroxide: bridge for protonmail, the smtp server

[Service]
ExecStart=hydroxide smtp

[Install]
WantedBy=default.target
