# .config/systemd/user/hydroxide-imap.service
[Unit]
Description=hydroxide: bridge for protonmail, the imap server

[Service]
ExecStart=hydroxide imap

[Install]
WantedBy=default.target
